# btrfs-blscfg.service
[Unit]
Description=BLSCFG generator for BTRFS snapshots

[Service]
Type=oneshot
ExecStart=/usr/local/bin/btrfs-blscfg
