# btrfs-blscfg.path
[Unit]
Description=BLSCFG generator for BTRFS snapshots

[Path]
PathChanged=/boot
PathChanged=/.snapshots

[Install]
WantedBy=default.target
